# services/meta/PiScan.service
[Unit]
Description=System script to run PiScan service

[Service]
Type=forking
ExecStart=/usr/bin/startup.sh

[Install]
WantedBy=multi-user.target
